# systemd/freeradius.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/freeradius -f -lstdout
